# config/crontab
*/2 * * * *     /usr/bin/ruby /home/mint/.bin/gridsquare.rb > /dev/null 2>&1
