# config/crontab
*/2 * * * *     /usr/bin/ruby /usr/bin/gridsquare.rb > /dev/null 2>&1
